# systemd/medical-ble-hub-backup.service
[Unit]
Description=Medical BLE Hub Database Backup

[Service]
Type=oneshot
User=sh3
Group=sh3
WorkingDirectory=/home/sh3/Desktop/sh3-hw-layer-experiments
ExecStart=/home/sh3/Desktop/sh3-hw-layer-experiments/hub_db_backup.sh
